# pi-setup/systemd/cnc-usb-mass-storage.service.tmpl
[Unit]
Description=CNC USB Mass Storage gadget (loads g_mass_storage backed by the share image)
DefaultDependencies=no
After=local-fs.target
Before=cnc-usb-watcher.service
ConditionPathExists=@@IMAGE_PATH@@

[Service]
Type=oneshot
RemainAfterExit=yes
# Stale-load guard: if the module is already loaded with a different file, unload first.
ExecStartPre=/bin/sh -c 'modprobe -r g_mass_storage 2>/dev/null || true'
#
# ro=0 — controller can write back to the stick (DPRNT logs, output files).
#
# This is safe under the v2 architecture because Linux NEVER mounts this
# image. The watcher detaches the LUN, syncs through mtools, and
# reattaches — so the host and the controller never touch the file at
# the same time. That's the documented-safe pattern from
# Documentation/usb/mass-storage.rst.
#
# (If you see corruption symptoms again, check that the watcher hasn't
# regressed to loop-mounting the image at $SHARE_PATH — `mountpoint -q
# $SHARE_PATH` should return false.)
#
ExecStart=/sbin/modprobe g_mass_storage \
  file=@@IMAGE_PATH@@ \
  removable=1 \
  ro=0 \
  stall=0 \
  cdrom=0 \
  nofua=1 \
  iManufacturer="@@USB_VENDOR@@" \
  iProduct="@@USB_PRODUCT@@" \
  iSerialNumber="@@USB_SERIAL@@"
ExecStop=/sbin/modprobe -r g_mass_storage

[Install]
WantedBy=multi-user.target
