# pi-setup/systemd/filebrowser.service.tmpl
[Unit]
Description=File Browser (NC fork) — web file manager pointed at the CNC share
After=network-online.target local-fs.target
Wants=network-online.target

[Service]
Type=simple
User=@@FB_USER@@
Group=@@FB_USER@@
WorkingDirectory=@@FB_WORKDIR@@
ExecStart=@@FB_BIN@@ --root @@SHARE_PATH@@ --address 0.0.0.0 --port 8080 --database @@FB_DB@@
Restart=on-failure
RestartSec=3

# Hardening — light. Filebrowser's whole job is to serve a folder
# (typically under /home), and it persists state in ~/.config, so
# locking down /home blocks both serving and the DB write.
NoNewPrivileges=yes
ProtectSystem=full
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
