# docker/root/etc/services.d/filebrowser/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid abc filebrowser -c /config/settings.json -d /database/filebrowser.db;